# auto-updater/DEBIAN/postinst
#!/bin/bash
systemctl daemon-reexec
systemctl daemon-reload
systemctl enable auto-updater.service
