# auto-updater/DEBIAN/postrm
#!/bin/bash

systemctl disable auto-updater.service
rm -f /etc/systemd/system/auto-updater.service
systemctl daemon-reload
